# server/deploy/jianghu.service
[Unit]
Description=Jianghu Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/jianghuxiaobao/data
ExecStart=/usr/local/bin/jianghu-server
Restart=on-failure
RestartSec=5s

# 配置文件路径（相对于 WorkingDirectory）
Environment="CONFIG_PATH=/root/jianghu/config.yaml"

[Install]
WantedBy=multi-user.target
